# certbot/etc/systemd/system/renew.timer
[Unit]
Description=Timer for the certificate renewal service

[Timer]
OnCalendar=Weekly
Persistent=true

[Install]
WantedBy=timers.target