# certbot/etc/systemd/system/renew.service
[Unit]
Description=Проверка и обновление сертификата при помощи certbot

[Service]
Type=simple
ExecStart=/usr/local/bin/renew

[Install]
WantedBy=multi-user.target